# plugins/smartcard/gsd-smartcard.service.in
[Unit]
Description=Smartcard handling
PartOf=graphical-session.target
After=gnome-session.service graphical-session-pre.service

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-smartcard
Restart=on-failure
BusName=@plugin_dbus_name@

[Install]
WantedBy=gnome-settings-daemon.target
