# plugins/screensaver-proxy/gsd-screensaver-proxy.service.in
[Unit]
Description=Freedesktop screensaver handling
PartOf=graphical-session.target
After=@default_after@

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-screensaver-proxy
Restart=on-failure
BusName=@plugin_dbus_name@

[Install]
WantedBy=gnome-settings-daemon.target
