# plugins/screensaver-proxy/gsd-screensaver-proxy.service.in
[Unit]
Description=Freedesktop screensaver handling
PartOf=graphical-session.target
After=gnome-session.service graphical-session-pre.service

[Service]
Type=notify
NotifyAccess=main
ExecStart=@libexecdir@/gsd-screensaver-proxy
Restart=on-failure

[Install]
WantedBy=gnome-settings-daemon.target
