# plugins/color/gsd-color.service.in
[Unit]
Description=Color management
PartOf=graphical-session.target
After=@default_after@

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-color
Restart=on-failure
BusName=@plugin_dbus_name@

[Install]
WantedBy=gnome-settings-daemon.target gnome-settings-daemon-gdm.target
