# plugins/color/gsd-color.target.in
[Unit]
Description=GNOME Color management

BindsTo=gsd-color.service
After=gsd-color.service

Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target
