# plugins/sharing/gsd-sharing.service.in
[Unit]
Description=Sharing handling
PartOf=graphical-session.target
After=gnome-session-bus.target graphical-session-pre.target gnome-shell-wayland.target

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-sharing
Restart=on-failure
BusName=@plugin_dbus_name@

[Install]
WantedBy=gnome-settings-daemon.target
