# plugins/datetime/gsd-datetime.service.in
[Unit]
Description=Date & Time handling
PartOf=graphical-session.target
After=gnome-session.service graphical-session-pre.service

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-datetime
Restart=on-failure
BusName=@plugin_dbus_name@

[Install]
WantedBy=gnome-settings-daemon.target
