# plugins/keyboard/gsd-keyboard.target.in
[Unit]
Description=GNOME Keyboard handling

BindsTo=gsd-keyboard.service
After=gsd-keyboard.service

Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target
