# plugins/gsd.target.in
[Unit]
Description=@description@ target
CollectMode=inactive-or-failed

# Pull in the service
Wants=gsd-@plugin_name@.service

# Require GNOME session and specify startup ordering
Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target
