# plugins/gsd.service.in
[Unit]
Description=@description@ service
CollectMode=inactive-or-failed
RefuseManualStart=true
RefuseManualStop=true

After=gnome-session-initialized.target

# Requisite/PartOf means the dependency is not loaded automatically.
# The ordering here also implies Before=gnome-session.target
Requisite=@plugin_dbus_name@.target
PartOf=@plugin_dbus_name@.target
Before=@plugin_dbus_name@.target

@plugin_gate_units_section@

[Service]
Slice=session.slice
Type=dbus
ExecStart=@libexecdir@/gsd-@plugin_name@
Restart=on-failure
BusName=@plugin_dbus_name@
TimeoutStopSec=5
# We cannot use OnFailure as e.g. dependency failures are normal
# https://github.com/systemd/systemd/issues/12352
ExecStopPost=@libexecdir@/gnome-session-ctl --exec-stop-check
