# plugins/wacom/gsd-wacom.target.in
[Unit]
Description=GNOME Wacom handling

BindsTo=gsd-wacom.service

Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target

# Use the following once wacom.target or similar exists
#Requisite=wacom.target
#After=wacom.target
#PartOf=wacom.target

[Install]
WantedBy=gnome-session-initialized.target
#WantedBy=wacom.target
