# plugins/media-keys/gsd-media-keys.target.in
[Unit]
Description=GNOME Media keys handling

BindsTo=gsd-media-keys.service
After=gsd-media-keys.service

Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target
