# plugins/media-keys/gsd-media-keys.service.in
[Unit]
Description=Media keys handling
PartOf=graphical-session.target
After=gnome-session-bus.target graphical-session-pre.target gnome-shell-wayland.target pulseaudio.service
Requires=pulseaudio.service

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-media-keys
Restart=on-failure
BusName=@plugin_dbus_name@

[Install]
WantedBy=gnome-settings-daemon.target gnome-settings-daemon-gdm.target
