# plugins/sound/gsd-sound.target.in
[Unit]
Description=GNOME Sound sample caching handling

BindsTo=gsd-sound.service
After=gsd-sound.service

Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target

Requires=pulseaudio.service
After=pulseaudio.service
# Keep running when pulseaudio.service stops
