# plugins/clipboard/gsd-clipboard.target.in
[Unit]
Description=GNOME Clipboard handler

BindsTo=gsd-clipboard.service

Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target
