# plugins/power/gsd-power.service.in
[Unit]
Description=GNOME Power management handling
# Note that we do the shengians with having a target and a service solely
# so that OnFailure is not called in case of a dependency failure.
OnFailure=gnome-session-failed.target
RefuseManualStart=true
RefuseManualStop=true

PartOf=gsd-power.target

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-power
Restart=on-failure
BusName=@plugin_dbus_name@
