# plugins/dummy/gsd-dummy.service.in
[Unit]
Description=GNOME Settings Daemon dummy replacement plugin (@pluginname@)
PartOf=graphical-session.target
After=gnome-session-bus.target graphical-session-pre.target gnome-shell-wayland.target

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-dummy --dummy-name=@pluginname@
Restart=on-failure
BusName=@plugin_dbus_name@

[Install]
WantedBy=gnome-settings-daemon.target
