# plugins/dummy/gsd-dummy.service.in
[Unit]
Description=GNOME Settings Daemon dummy replacement plugin (@pluginname@)
PartOf=graphical-session.target
After=gnome-session.service graphical-session-pre.service

[Service]
Type=notify
NotifyAccess=main
ExecStart=@libexecdir@/gsd-dummy --dummy-name=@pluginname@
Restart=on-failure

[Install]
WantedBy=gnome-settings-daemon.target
