# plugins/dummy/gsd-dummy.target.in
[Unit]
Description=GNOME Settings Daemon dummy replacement plugin (@pluginname@)

BindsTo=gsd-dummy.service
After=gsd-dummy.service

Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target

# See require-started.target on how to require another unit to be already running
# See dependency.target on how to start a dependency and require it to be present
#
#Requisite=require-started.target
#After=require-started.target
#PartOf=require-started.target
#Requires=dependency.target
#After=dependency.target
#
# Also add PartOf if this unit should be stopped together with dependency.target
#PartOf=dependency.target

# This is what the install section would look like, but we create the symlinks
# through meson
#[Install]
#WantedBy=gnome-session-initialized.target
##WantedBy=require-started.target
