# data/gnome-settings-daemon-gdm.target
[Unit]
Description=GNOME Settings daemon (GDM)
StopWhenUnneeded=yes
Conflicts=gnome-settings-daemon.target
