# data/gnome-settings-daemon.target
[Unit]
Description=GNOME Settings daemon
StopWhenUnneeded=yes
RefuseManualStart=yes
