# plugins/print-notifications/gsd-print-notifications.target.in
[Unit]
Description=GNOME Printer notifications

Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-print-notifications
Restart=on-failure
BusName=@plugin_dbus_name@

# We cannot use OnFailure as that may be called for dependency issues
ExecStopPost=@libexecdir@/gnome-session-ctl --fail-if-start-limit-hit

[Install]
WantedBy=gnome-session-initialized.target
