# plugins/print-notifications/gsd-print-notifications.service.in
[Unit]
Description=Print dbus
PartOf=graphical-session.target
After=@default_after@

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-print-notifications
Restart=on-failure
BusName=@plugin_dbus_name@

[Install]
WantedBy=gnome-settings-daemon.target
