# plugins/wwan/gsd-wwan.target.in
[Unit]
Description=GNOME WWan management

BindsTo=gsd-wwan.service
After=gsd-wwan.service

Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target
