# plugins/xsettings/gsd-xsettings.service.in
[Unit]
Description=XSettings
PartOf=graphical-session.target
After=gnome-session-bus.target graphical-session-pre.target gnome-shell-wayland.target

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-xsettings
Restart=on-failure
BusName=@plugin_dbus_name@

[Install]
WantedBy=gnome-settings-daemon.target gnome-settings-daemon-gdm.target
