# plugins/xsettings/gsd-xsettings.target.in
[Unit]
Description=GNOME XSettings

BindsTo=gsd-xsettings.service
After=gsd-xsettings.service

Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target

Requisite=gnome-session-x11-services.target
After=gnome-session-x11-services.target
PartOf=gnome-session-x11-services.target
