# plugins/rfkill/gsd-rfkill.service.in
[Unit]
Description=RFKill handling
PartOf=graphical-session.target
After=@default_after@

[Service]
Type=dbus
ExecStart=@libexecdir@/gsd-rfkill
Restart=on-failure
BusName=@plugin_dbus_name@

[Install]
WantedBy=gnome-settings-daemon.target
